# ops/odyssee-listener.service
[Unit]
Description=Firefly (Odyssee) Listener
After=network.target

[Service]
Type=simple
User=odyssee
WorkingDirectory=/opt/odyssee
EnvironmentFile=/etc/odyssee/.env
ExecStart=/usr/bin/node /opt/odyssee/dist/listener.js
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=firefly-listener

[Install]
WantedBy=multi-user.target